# kreep.service
[Unit]
Description=Kreep - Self-hosted note taking app
After=network.target

[Service]
Type=simple
User=your-username
WorkingDirectory=/path/to/kreep
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=kreep

Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
