# kreep.service
[Unit]
Description=Keep Clone - Self-hosted note taking app
After=network.target

[Service]
Type=simple
User=your-username
WorkingDirectory=/path/to/kreep
# Rebuild the fingerprinted frontend bundles before starting, so a git pull that
# changed public/app.js is actually served instead of a stale public/dist.
ExecStartPre=/usr/bin/npm run build
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kreep

Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
